# local_file/var/lib/dpkg/info/php7.2-fpm.conffiles
/etc/apache2/conf-available/php7.2-fpm.conf
/etc/init.d/php7.2-fpm
/etc/init/php7.2-fpm.conf
/etc/logrotate.d/php7.2-fpm
/etc/php/7.2/fpm/php-fpm.conf
/etc/php/7.2/fpm/pool.d/www.conf
